# packages/httpd/source/system.d/service.mysql-server.service
[Unit]
Description=MySQL server
After=network.target

[Service]
Type=forking
PIDFile=/var/logs/mysqld.pid
ExecStart=/bin/sh -c "exec sh /storage/.xbmc/addons/service.web.httpd/bin/mysql-server.start"
ExecStop=/bin/sh -c "exec sh /storage/.xbmc/addons/service.web.httpd/bin/mysql-server.stop"
TimeoutStopSec=5
Restart=always
RestartSec=10
StartLimitInterval=0

[Install]
WantedBy=xbmc.target
