# packages/iMX6/status-led/system.d/status-led.service
[Unit]
Description=Status LED
After=xbmc.service

[Service]
Type=oneshot
ExecStart=-/bin/sh -c "/usr/bin/status-led on"
ExecStop=-/bin/sh -c "/usr/bin/status-led blinking"
RemainAfterExit=yes

[Install]
WantedBy=xbmc.target
