# packages/iMX6/soc-fan/system.d/soc-fan-monitor.service
[Unit]
Description=SOC fan monitor service

[Service]
ExecStart=/usr/bin/soc-fan-monitor

[Install]
WantedBy=multi-user.target
